# last-hub.service
[Unit]
Description=LAST Hub - Fastify Service
After=network.target

[Service]
Type=simple
User=ocs
WorkingDirectory=%h/last-hub
ExecStart=/usr/bin/node %h/last-hub/index.js
Restart=on-failure
Environment=http_proxy=http://bcproxy.weizmann.ac.il:8080
Environment=https_proxy=http://bcproxy.weizmann.ac.il:8080
Environment=no_proxy=localhost,127.0.0.1
Environment=NODE_ENV=production
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
